# puppet/modules/acquia_base/templates/docker/docker-storage-setup.service.erb
[Unit]
Description=Docker Storage Setup
<% @aws_block_devices.each do |device|
     block_device = device.gsub(%r{^/dev/}, '')
-%>
After=device-setup-<%= block_device %>.service
<% end -%>
Before=docker.service

[Service]
Type=oneshot
ExecStart=/usr/bin/docker-storage-setup
EnvironmentFile=-/etc/sysconfig/docker-storage-setup

[Install]
WantedBy=multi-user.target